# etc/fastishcgi.service
[Unit]
Description = FastCGI <-> CGI bridge 

[Service]
ExecStart = /usr/bin/fastishcgid

[Install]
WantedBy = multi-user.target

